# asn_server/Testbench/system/asnnode63/rpiserver.service
[Unit]
Description=Raspberry Pi Sound Server

[Service]
Type=simple
ExecStart=/usr/bin/rpiserver

[Install]
WantedBy=multi-user.target 


